# scripts/mashoodwear-api.service.example
[Unit]
Description=Mashoodwear API (Express)
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mashoodwear/backend
EnvironmentFile=/opt/mashoodwear/backend/.env
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
